# install/custom_services/notification_service.timer
[Unit]
Description=Run script on boot

[Timer]
OnCalendar=*-*-* 4:00:00
Persistent=true

[Install]
WantedBy=timers.target
